# packages/devnet/files/etc/uci-defaults/90_custom-fw-repo
#!/bin/sh
uci set eupgrade.main.custom_api_url='http://10.13.0.2:8081/'
uci commit eupgrade
exit 0
